# distcc-configurator/distcc-configurator
#!/sbin/runscript

RC_REBOOT="NO"
ADDRESS=$(ifconfig | grep netmask | sed -ne '1p' | awk '{print $2}')
NETMASK=$(ifconfig | grep netmask | sed -ne '1p' | awk '{print $4}')
NETWORK=$(ipcalc ${ADDRESS}/${NETMASK} | grep Network | awk '{print $2}')

depend() {
  need net
}

start() {
	firstrun
	if ! [ ${RC_CMD} == "restart" ] || ! [ ${RC_REBOOT} == "YES" ];then
		ebegin "configuring distcc hosts"
		configure
	fi
	ebegin "starting distcc-config-listener"
		start-stop-daemon --start --make-pidfile --exec /usr/bin/listener --background \
			--pidfile /etc/distcc/listener.pid
	eend $?
	/etc/init.d/distccd start
}

stop() {
	if [ ${RC_RUNLEVEL} == "shutdown" ] && ! [ ${RC_REBOOT} == "YES" ];then rm /etc/distcc/lock;fi
	ebegin "stopping distcc-configurator"
  if ! [ ${RC_CMD} == "restart" ] || ! [ ${RC_REBOOT} == "YES" ];then
		configure
	fi
    start-stop-daemon --stop --exec /usr/bin/listener \
      --pidfile /etc/distcc/listener.pid
  eend $?
  /etc/init.d/distccd stop
}

firstrun() {
	if ! [ -f /etc/distcc/lock ];then
		touch /etc/distcc/lock
		backup_config
		write_config
		add_netrange
	fi
}

bawler() {
	for i in $(seq 0 $((${#HOSTS[*]}-1)));do
		if ! [ ${ADDRESS} == ${HOSTS[i]} ];then
			echo "${NEW_HOSTS}" | nc -w 1 ${HOSTS[i]} 3633
		fi
	done
}

write_config() {
cat > /etc/conf.d/distccd << 'EOF'
NETWORK=
DISTCCD_OPTS=""
DISTCCD_EXEC="/usr/bin/distccd"
DISTCCD_PIDFILE="/var/run/distccd/distccd.pid"
DISTCCD_OPTS="${DISTCCD_OPTS} --port 3632"
DISTCCD_OPTS="${DISTCCD_OPTS} --log-level critical"
DISTCCD_OPTS="${DISTCCD_OPTS} --allow ${NETWORK}"
DISTCCD_OPTS="${DISTCCD_OPTS} -N 15"
EOF
}

add_netrange() {
	sed -i '/NETWORK=/d' /etc/conf.d/distccd
	sed -i "1 i\NETWORK=${NETWORK}" /etc/conf.d/distccd
}

calculate_hosts() {
	HOSTS_NEW=$(nmap -oG - -sT -p 3632 ${NETWORK} |grep open | awk '{print $2 " "}' | tr -d \\012)
	if [ ${RC_CMD} == "start" ];then
		NEW_HOSTS="${ADDRESS} ${HOSTS_NEW}"
	elif [ ${RC_CMD} == "stop" ];then
		NEW_HOSTS=`echo ${HOSTS_NEW} | sed "s/${ADDRESS}//g" | sed 's/^[ \t]*//;s/[ \t]*$//'`
	fi
	HOSTS=(${NEW_HOSTS})
}

backup_config() {
		grep "# /etc/conf.d/distccd: config file for /etc/init.d/distccd" /etc/conf.d/distccd 1> /dev/null
		if [ $? -eq 0 ];then mv /etc/conf.d/distccd /etc/conf.d/distccd.gentoo.org;fi
}

configure() {
	if ! [ ${RC_CMD} == "restart" ] || ! [ ${RC_REBOOT} == "YES" ];then
		calculate_hosts
		bawler
	fi
}
